# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn ahps_api:app -b 0.0.0.0:5000
directory=/src